# crontabs/crontab-pi
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
# Avoid sending mail
MAILTO=""

# panorama 
3 * * * *    cd /home/pi/panogen && python3 panogen.py | logger -t panogen
#webshot
17 * * * *    cd /home/pi/weather_station/webshot && ./run.sh | logger -t webshot
