# crontabs/crontab-root
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
# Avoid sending mail
MAILTO=""

# webcam: generate website images from raw webcam shots, and upload them to ftp
*/5 * * * *    cd /home/pi/weather_station/webcam/scripts && chronic ./gen_new_pic.sh | logger -t webcam
